# .config/profile.d/plugins.sh
# >>> fzf integration >>>
source <(fzf --bash)
# <<< fzf integration <<<

# >>> direnv integration >>>
eval "$(direnv hook bash)"
# <<< direnv integration <<<
